# platform/docker_images/vpnsecret/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D -e

[program:vpnsecret]
command=python3 /server/vpnsecret.py

# [program:krillpubd]
# command=/usr/local/bin/krillpubd -c /var/krill/krillpubd.conf
