# platform/docker_images/irrd/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D -e

[program:irrd]
command=/bin/bash -c "/usr/bin/irrd_database_upgrade && /usr/bin/irrd"
