# platform/docker_images/routinator/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D -e

[program:routinator]
command=/usr/local/bin/routinator server
