# platform/docker_images/switch/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D -e

[program:ovs]
command=/usr/local/bin/run_ovs
