# platform/docker_images/router/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D -e

[program:frr]
command=/usr/local/bin/run_frr

[program:looking_glass]
command=/usr/local/bin/looking_glass
